# systemd/nft-watchdog.service
[Unit]
Description=NFT Firewall Watchdog — WireGuard VPN health monitor
After=network-online.target wg-quick@wg0.service
Wants=network-online.target wg-quick@wg0.service

[Service]
Type=simple
User=fw-admin
WorkingDirectory=/opt/nft-firewall
# Attempt to start VPN if config exists (requires sudoers grant for fw-admin)
ExecStartPre=-/usr/bin/sudo /usr/bin/systemctl start wg-quick@wg0.service
ExecStart=/usr/bin/python3 /opt/nft-firewall/src/main.py watchdog daemon
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nft-watchdog

[Install]
WantedBy=multi-user.target
