# systemd/nft-threat-update.service
[Unit]
Description=NFT Firewall — Threat Feed Sync
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=fw-admin
UMask=0077
PrivateTmp=true
ProtectClock=true
ProtectControlGroups=true
ProtectHostname=true
ProtectKernelLogs=true
ProtectKernelModules=true
LockPersonality=true
RestrictRealtime=true
SystemCallArchitectures=native
WorkingDirectory=/opt/nft-firewall
ExecStart=/usr/bin/sudo -n /usr/local/bin/fw threat-update
Environment=PYTHONPATH=/opt/nft-firewall/src
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nft-threat-update
