# systemd/nft-firewall-threatfeed.service
[Unit]
Description=NFT Firewall threat feed update

[Service]
Type=oneshot
User=fw-admin
ExecStart=/usr/local/bin/fw threat-update
