# systemd/nft-threat-update.service
[Unit]
Description=NFT Firewall — Threat Feed Sync
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=fw-admin
WorkingDirectory=/opt/nft-firewall
ExecStart=/usr/bin/python3 /opt/nft-firewall/src/main.py threat-update
Environment=PYTHONPATH=/opt/nft-firewall/src
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nft-threat-update
